# day1/task1/99-log-USB-drive.rules
ACTION=="add", ATTRS{serial}=="558582C6333542CB61AE778B",RUN+="/usr/bin/trigger.sh"
